# udev/em0018_slot7.rules
ACTION=="add|remove", SUBSYSTEM=="gpio", DEVPATH=="*/i2c-1/1-0027/gpiochip*",  RUN+="/opt/unipi/tools/uhelper RO_EMO7 gpiochip 7 $name"